# scripts/eval-baseline-launchd.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  scripts/eval-baseline-launchd.plist.template — SMI-4764 Wave 2

  macOS LaunchAgent template for the canonical-dev retrieval-eval cron.

  Setup (canonical dev only — see .claude/development/eval-cron-setup.md):

    1. Copy template, substituting __REPO_PATH__, __HOME__, and
       __PROJECT_DIR_ENCODED__ (SKILLSMITH_PROJECT_DIR_ENCODED, read from .env via
       varlock AT SETUP TIME — the cron itself does NOT run under varlock, which
       hangs in launchd's keychain-less context):
         PROJ=$(varlock run -- sh -c 'printf %s "$SKILLSMITH_PROJECT_DIR_ENCODED"')
         sed "s|__REPO_PATH__|$(pwd)|g; s|__HOME__|$HOME|g; s|__PROJECT_DIR_ENCODED__|$PROJ|g" \
           scripts/eval-baseline-launchd.plist.template \
           > ~/Library/LaunchAgents/app.skillsmith.eval-baseline-cron.plist

    2. Load:
         launchctl unload ~/Library/LaunchAgents/app.skillsmith.eval-baseline-cron.plist 2>/dev/null
         launchctl load   ~/Library/LaunchAgents/app.skillsmith.eval-baseline-cron.plist

    3. Verify:
         launchctl list | grep skillsmith
         tail ~/.skillsmith/logs/eval-cron-*.log

  Cadence: weekly Sunday 03:00 local time (matches plan §3 — quiet hours,
  outside of Mon 06:00 UTC eval-weekly.yml so audit/cron don't collide).

  Notes:
    - StandardOutPath / StandardErrorPath capture launchd stdout (the
      script itself writes detailed logs to ~/.skillsmith/logs/).
    - RunAtLoad=false: do NOT fire on `launchctl load`. Schedule-only.
    - StartCalendarInterval is the launchd equivalent of cron expressions.
    - SMI-5353: launchd has no per-invocation timeout key equivalent to
      systemd's TimeoutStartSec. If a run wedges (hung `docker compose run`
      or stuck eval), kill it manually:
        launchctl list | grep skillsmith         # find the PID
        launchctl kill TERM gui/$(id -u)/app.skillsmith.eval-baseline-cron
      The script's EXIT trap still tears the ephemeral eval container down on
      kill. As a backstop, `launchctl remove app.skillsmith.eval-baseline-cron`
      unloads the agent entirely.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>app.skillsmith.eval-baseline-cron</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <!-- SMI-5353: the decoupled cron starts its OWN container and needs
         SKILLSMITH_PROJECT_DIR_ENCODED (for the /skillsmith-memory bind-mount) —
         provided via EnvironmentVariables below, NOT by wrapping in `varlock`.
         varlock HANGS in launchd's keychain-less background context (it never
         exec's the cron). gh/git/docker work in launchd natively, exactly as the
         pre-SMI-5353 cron did. -->
    <string>cd __REPO_PATH__ &amp;&amp; ./scripts/eval-baseline-cron.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_PATH__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__HOME__/.skillsmith/logs/eval-cron-launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.skillsmith/logs/eval-cron-launchd-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- PATH for docker/gh/git — launchd's `bash -l` env is minimal. Add explicit
         entries here if the canonical dev's tools live elsewhere. -->
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- SMI-5353: the eval container's /skillsmith-memory bind-mount needs this
         (a path, NOT a secret). Substituted at setup from .env via varlock. -->
    <key>SKILLSMITH_PROJECT_DIR_ENCODED</key>
    <string>__PROJECT_DIR_ENCODED__</string>
  </dict>
</dict>
</plist>
